# src/network_api/services/network_api.service
[Unit]
Description=Network API Service
After=network.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/network_api_node
Restart=always
User=root
# Could set CPU affinity or real-time priority if needed

[Install]
WantedBy=multi-user.target
